# conf/crontab
*/2 * * * * cd /root/vmlist && bin/make_vmlist.rb